# sources/meta-oes-service/recipes-brightness/brightness/files/brightness.service
[Unit]
Description=Brightness monitor service

[Service]
ExecStart=/etc/brightnessControl.sh

[Install]
WantedBy=multi-user.target

