# sources/meta-oes-service/recipes-brightness/brightness/files/ip_check.service
[Unit]
Description=ip address monitor service

[Service]
ExecStart=/etc/ip_check.sh

[Install]
WantedBy=multi-user.target

