# sources/meta-oes-service/recipes-ipcheck/ipcheck/files/ip_check.service
[Unit]
Description=ip address monitor service
After=startup.service

[Service]
ExecStart=/etc/ip_check.sh

[Install]
WantedBy=multi-user.target

