# sources/meta-oes/recipes-startupservice/startupservice/files/startup.service
[Unit]
Description=Start up service

[Service]
ExecStart=/etc/oesStart.sh

[Install]
WantedBy=multi-user.target
